# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronisation timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m
Unit=mbsync.service

[Install]
WantedBy=timers.target
